# etc/profile.d/desktop-misc.sh
# mesa 25.2-rc1 broke vulkan rendering in GTK4/libadwaita
# so let's use opengl for now
export GSK_RENDERER="gl"
