# etc/profile.d/makemkv-libmmdb.sh
# Use MakeMKV's libmmbd as default for disc decryption
export LIBAACS_PATH="libmmbd"
export LIBBDPLUS_PATH="libmmbd"
