# etc/profile.d/plasma-preload.sh
# "preload popups of every plasmoid"
# https://community.kde.org/Plasma/Environment_Variables
export PLASMA_PRELOAD_POLICY="aggressive"
