# etc/profile.d/vulkan-video.sh
# enable experimental vulkan video
export ANV_VIDEO_ENCODE=1
export ANV_VIDEO_DECODE=1 
