# etc/profile.d/vulkan-video.sh
# enable experimental vulkan video
export ANV_DEBUG="${ANV_DEBUG+${ANV_DEBUG},}video-decode,video-encode"
